# etc/supervisor/conf.d/worker.conf
[program:worker]
command=php /app/bin/console app:run
autorestart=true
autostart=true
stdout_logfile=/var/log/soutput.log
stderr_logfile=/var/log/soutput.log
